# fleet-drop-in.service
[Unit]
After=aws-fleet-metadata.service

[Service]
EnvironmentFile=/etc/fleet-metadata
